# pic_programmer_v1/dist/pic_programmer.conf
description "PIC Programming Service"
author "Marco Antonio Garcia"

start on filesystem or runlevel [2345]
stop on shutdown

script
    echo $$ >> /var/run/pic_programmer.pid
    exec /home/g2marco/dist/pic_programmer
end script

pre-start script
    echo "[`date`] pic_programmer starting" >> /var/log/pic_programmer.log
end script

pre-stop script
    rm /var/run/pic_programmer.pid
    echo "[`date`] pic_programmer stopping" >> /var/log/pic_programmer.log
end script

